# jnwebsphere/files/default/init/ibm_websphere_node
#! /bin/sh
###
# chkconfig: 235 98 55
# description: Manages the services you are controlling with the chkconfig command
###


case "$1" in
  start)
          /opt/IBM/WebSphere/AppServer/profiles/Node01/bin/startNode.sh
          ;;
  stop)
          /opt/IBM/WebSphere/AppServer/profiles/Node01/bin/stopNode.sh
          ;;
  restart)
          $0 stop
          sleep 2
          $0 start
          ;;
  *)
        echo "Usage: /sbin/service ibm_websphere_node {start|stop|restart}"
        exit 1
esac


exit 0
